# packaging/deb/DEBIAN/postinst
#!/bin/sh
set -e

# Create system user
if ! id -u crosstalk >/dev/null 2>&1; then
  useradd --system --no-create-home --shell /usr/sbin/nologin \
    --home-dir /var/lib/crosstalk crosstalk
fi

# Create base directories (transport is cloned by `crosstalk install <git-url>`)
mkdir -p /etc/crosstalk /var/lib/crosstalk/workspaces /var/lib/crosstalk/.ssh
chown -R crosstalk /var/lib/crosstalk
chmod 700 /var/lib/crosstalk/.ssh

# Enable service unit (don't start — no transport yet)
systemctl daemon-reload
systemctl enable crosstalk.service >/dev/null 2>&1 || true

echo ""
echo "crosstalk-runtime installed."
echo "Next: run the setup command with your transport repo URL:"
echo "  sudo crosstalk install <git-url>"
echo ""
