# packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

systemctl stop crosstalk.service >/dev/null 2>&1 || true
systemctl disable crosstalk.service >/dev/null 2>&1 || true
systemctl daemon-reload || true
